# application-integration/container/eks/amazon-eks-autonomous-driving-data-service/scripts/webviz.service
[Unit]
Description=webviz service

[Service]
User=ubuntu
ExecStart=/usr/local/bin/start-webviz.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
